# templates/mailhog/mailhog.service.j2
[Unit]
Description=Mailhog Daemon
After=network.target

[Service]
ExecStart=/usr/local/bin/MailHog
StandardOutput=journal

[Install]
WantedBy=multi-user.target
