# templates/start/upstart.conf.j2
# CI/CD managed
# Do NOT modify this file by hand!

description "{{item.description|default('The service')}}"

start on runlevel [2345]
stop on runlevel [016]

{% for key in item.env|default({}) %}
env {{key}}={{item.env[key]}}
{% endfor %}

chdir {{item.workingdirectory|default('/tmp')}}
setuid {{item.user|default('root')}}
exec {{item.execstart}}
{% for line in item.extralines|default([]) %}
{{ line|default('') }}
{% endfor %}
respawn
console log
