# templates/start/upstart.conf.j2
# CI/CD managed
# Do NOT modify this file by hand!
# Uses the same conf. structure as systemd.conf.j2, automatic env. variable PIDFile

description "{{item.startup.description|default('The service')}}"

start on runlevel [2345]
stop on runlevel [016]

{% for key in item.startup.env|default({}) %}
env {{key}}={{item.startup.env[key]}}
{% endfor %}
env PIDFile={{item.startup.pidfile|default('/tmp/' + item.startup.name + '.pid')}}
env WORKDIR={{item.startup.workingdirectory|default('/tmp')}}

chdir {{item.startup.workingdirectory|default('/tmp')}}
setuid {{item.startup.user|default('root')}}
exec {{item.startup.execstart}}
{% if item.startup.execstop|default(False) %}
pre-stop exec {{item.startup.execstop}}
{% endif %}
{% for line in item.startup.extralines|default([]) %}
{{ line|default('') }}
{% endfor %}
respawn
console log
